# suspend_lock/lock_session-suspend.service
[Unit]
Before=suspend.target hybrid-sleep.target

[Service]
Type=oneshot
ExecStart=/usr/bin/lock-session lock
