# msmtp-systemd-path/msmtp.service
[Service]
Type=oneshot
ExecStart=msmtp-flush-mail %h/.local/share/msmtp-queue
Restart=on-failure
RestartSec=10min
