# msmtp-systemd-path/msmtp.path
[Path]
DirectoryNotEmpty=%h/.local/share/msmtp-queue
MakeDirectory=yes
DirectoryMode=700

[Install]
WantedBy=paths.target
